# platforms/organelle_m/rootfs/home/music/.config/systemd/user/start-mother.service
[Unit]
Description=Organelle M Mother

[Service]
ExecStart=/home/music/fw_dir/scripts/start-mother.sh

[Install]
WantedBy=default.target
