# platforms/organelle_m/rootfs/home/music/.config/systemd/user/status-led.service
[Unit]
Description=Disable status LED

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo none > /sys/class/leds/led0/trigger"

[Install]
WantedBy=default.target
